# services/aeyeweb-certbot-renew.timer
[Unit]
Description=Twice-daily Let's Encrypt renewal check for A-eye

[Timer]
OnCalendar=*-*-* 03,15:00:00
Persistent=true
Unit=aeyeweb-certbot-renew.service

[Install]
WantedBy=timers.target
